# lib/proxy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:register-dns]
command=/app/node_modules/.bin/tsx register-dns.ts
directory=/app
autostart=true
autorestart=false
startsecs=0
priority=1
stderr_logfile=/var/log/supervisor/register-dns.err.log
stdout_logfile=/var/log/supervisor/register-dns.out.log

[program:nginx]
command=/app/node_modules/.bin/tsx start-nginx.ts
directory=/app
autostart=true
autorestart=true
priority=10
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:health-check]
command=/app/node_modules/.bin/tsx health-check.ts
directory=/app
autostart=true
autorestart=false
priority=20
stderr_logfile=/var/log/supervisor/health-check.err.log
stdout_logfile=/var/log/supervisor/health-check.out.log

