# lib/proxy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/app/node_modules/.bin/tsx start-nginx.ts
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log

[program:health-check]
command=/app/node_modules/.bin/tsx health-check.ts
directory=/app
autostart=true
autorestart=false
stderr_logfile=/var/log/supervisor/health-check.err.log
stdout_logfile=/var/log/supervisor/health-check.out.log

